# debug-env/systemd-networkd.service.d/20-exec-start.conf
[Service]
ExecStart=
ExecStart=/opt/systemd-networkd-ubuntu/build/systemd-networkd
